# debian/postinst
#!/bin/bash
set -e

configure() {
    systemctl --user --global enable headphone-manager
}

case "$1" in
	"configure")
		configure
		;;
esac
