# debian/prerm
#!/bin/bash
set -e

remove() {
    systemctl --user --global disable headphone-manager
}

case "$1" in
	"remove")
		remove
		;;
esac
